# deploy/ansible/roles/services/templates/bot.conf
description "daemon for robocoup bot"

start on startup
stop on shutdown
respawn

NODE_ENV=production

script
  cd {{site_path}}
  /usr/bin/npm start
end script
