# deploy/ansible/roles/services/templates/bots.conf
description "daemon for bocoup bots"

start on startup
stop on shutdown
respawn

env NODE_ENV=production

script
  /usr/bin/npm start --prefix {{site_path}}
end script
