# deploy/ansible/roles/services/templates/bots.conf
description "daemon for bocoup bots"

start on startup
stop on shutdown
respawn

env NODE_ENV={{env}}
{% if env == "production" %}
{% for key in RDS.splitlines() %}
env {{key}}
{% endfor %}
{% for key in BOT_DB.splitlines() %}
env {{key}}
{% endfor %}
{% for key in BOT_KEYS.splitlines() %}
env {{key}}
{% endfor %}
{% for key in SES.splitlines() %}
env {{key}}
{% endfor %}
env RUN_JOBS=true
{% endif %}

chdir {{site_path}}

script
  /usr/bin/npm start
end script
