# deploy/ansible/roles/services/templates/bots.conf
description "daemon for bocoup bots"

start on startup
stop on shutdown
respawn

env NODE_ENV={{env}}
{% if env == "production" %}
env DB_HOST={{DB_HOST}}
env DB_USER={{DB_USER}}
env DB_NAME={{DB_NAME}}
env DB_PASSWORD={{DB_PASSWORD}}
env TOKEN_ROBOCUP={{TOKEN_ROBOCOUP}}
env TOKEN_THANKSBOT={{TOKEN_THANKSBOT}}
env EMAIL_SECRET={{EMAIL_SECRET}}
env EMAIL_KEY={{EMAIL_KEY}}
env RUN_JOBS=true
{% endif %}

script
  /usr/bin/npm start --prefix {{site_path}}
end script
